# sylph.service.example
[Unit]
Description=Sylph - Agent de recherche d'alternance
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/sylph
Environment="PYTHONPATH=/path/to/sylph"
# Utiliser le venv Python
ExecStart=/path/to/sylph/venv/bin/python3 /path/to/sylph/agent.py
Restart=always
RestartSec=10
StandardOutput=append:/path/to/sylph/logs/agent.log
StandardError=append:/path/to/sylph/logs/agent.log

[Install]
WantedBy=multi-user.target
